# kl.service
[Unit]
Description = kl

[Service]
ExecStart = /bin/kl

[Install]
WantedBy = multi-user.target
